# automation/com.linxz.whoisspending.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.linxz.whoisspending.daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/lxz/.local/share/whoIsSpendingDaily/repos/whoIsSpending/automation/daily_who_spending.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/lxz/.local/share/whoIsSpendingDaily/repos/whoIsSpending</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/lxz</string>
    <key>PATH</key>
    <string>/Users/lxz/.nvm/versions/node/v22.22.1/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TZ</key>
    <string>Asia/Shanghai</string>
    <key>LANG</key>
    <string>zh_CN.UTF-8</string>
    <key>LC_ALL</key>
    <string>zh_CN.UTF-8</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/lxz/Library/Logs/whoIsSpendingDaily/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/lxz/Library/Logs/whoIsSpendingDaily/launchd.err.log</string>
</dict>
</plist>
